# supervisord.conf
[supervisord]
nodaemon=true

[program:matt]
command=/matt/matt -c /matt/config/matt.toml

[program:nginx]
command=nginx -g "daemon off;"
